# user/claws@.service
[Unit]
Description = Claws

[Service]
ExecStart = /usr/bin/claws-mail
Environment=DISPLAY=%i

[Install]
WantedBy = mystuff.target