# user/cryptkeeper.service
[Unit]
Description = EncFS manager

[Service]
ExecStart = /usr/bin/cryptkeeper
Environment=DISPLAY=:0 

[Install]
WantedBy = mystuff.target