# user/e17.service
[Unit]
Description=your window manager service
Before=mystuff.target

[Service]
Environment=DISPLAY=:0
ExecStart=/usr/bin/enlightenment_start
Restart=always
RestartSec=10
 
[Install]
WantedBy=wm.target
