# user/get_iplayer.timer
[Unit]
Description=iplayer pvr

[Timer]
OnCalendar=06:30
Unit=get_iplayer.service

[Install]
WantedBy=default.target
