# user/get_iplayer.service
[Unit]
Description=get_iplayer PVR

[Service]
Type=oneshot
ExecStart=/usr/bin/get_iplayer --pvr

[Install]
WantedBy=mystuff.target