# user/hexchat@.service
[Unit]
Description = Hexchat

[Service]
ExecStart = /usr/bin/hexchat
Environment=DISPLAY=%i

[Install]
WantedBy = mystuff.target
