# user/mpd.socket
[Socket]
ListenStream=6600

[Install]
WantedBy=mystuff.target