# user/pidgin@.service
[Unit]
Description = Pidgin

[Service]
ExecStart = /usr/bin/pidgin
Environment=DISPLAY=%i

[Install]
WantedBy = mystuff.target
