# user/quasselcore.service
[Unit]
Description = quasselcore

[Service]
ExecStart = /usr/bin/quasselcore
Restart = always

[Install] 
WantedBy = mystuff.target