# user/smuxi-server.service
[Unit]
Description = smuxi-server

[Service]
ExecStart = /usr/bin/smuxi-server
Restart = always

[Install] 
WantedBy = mystuff.target